# pakkemaskin-skriver.service
[Unit]
Description=Pakkemaskin Skriver - pakkelapp-fangst
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pakkemaskin-skriver
# Juster --port og --baud. Bruk --bare-fangst første gang for å se lappformatet.
ExecStart=/usr/bin/python3 /home/pi/pakkemaskin-skriver/les_pakkelapp.py --port /dev/ttyUSB0 --baud 9600
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
