# python/no/pakkemaskin-skriver.service
[Unit]
Description=Pakkemaskin Skriver - pakkelapp-fangst
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rs232excel/python/no
# Juster --port og --baud. Bruk --bare-fangst første gang for å se lappformatet.
ExecStart=/usr/bin/python3 /home/pi/rs232excel/python/no/read_package.py --port /dev/ttyUSB0 --baud 9600 --usb-sti /media/usb0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
